# web-browser/systemd/agent-web-browser.service
[Unit]
Description=Agent Web Browser daemon

[Service]
ExecStart=/bin/sh -lc 'unit="$(readlink -f "$(systemctl --user show -P FragmentPath agent-web-browser.service)")"; root="$(dirname "$(dirname "$unit")")"; cd "$root" && exec /usr/bin/mise exec -- node daemon/browserd.js'
Environment=AGENT_WEB_CHROME_CDP_PORT=19339
Environment=AGENT_WEB_IDLE_TIMEOUT_MS=600000
NoNewPrivileges=yes
KillMode=control-group
TimeoutStopSec=10
